# rootfs/etc/services.d/dbus/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start D-Bus daemon (required for Avahi)
# ==============================================================================
bashio::log.info "Starting D-Bus daemon..."

# Create D-Bus directories
mkdir -p /var/run/dbus
mkdir -p /var/lib/dbus

# Initialize D-Bus machine ID if not exists
if [ ! -f /var/lib/dbus/machine-id ]; then
    dbus-uuidgen > /var/lib/dbus/machine-id
fi

# Clean up stale PID file
rm -f /var/run/dbus/pid

# Start D-Bus in foreground
exec dbus-daemon --system --nofork --nopidfile
