# rootfs/etc/services.d/cups-sharing/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Handle CUPS 3.0 Sharing Server shutdown
# ==============================================================================
bashio::log.info "CUPS 3.0 Sharing Server stopped"
