# rootfs/etc/services.d/cupsd/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Handle CUPS daemon shutdown
# ==============================================================================
bashio::log.info "CUPS daemon stopped"
