# rootfs/etc/services.d/cupsd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start CUPS 3 daemon with enhanced configuration
# ==============================================================================

bashio::log.info "Starting CUPS 3 (OpenPrinting) Enhanced..."

# Get all configuration values
admin_user=$(bashio::config 'admin_user')
admin_password=$(bashio::config 'admin_password')
log_level=$(bashio::config 'log_level')
share_printers=$(bashio::config 'share_printers')
browsing=$(bashio::config 'browsing')
allow_remote_admin=$(bashio::config 'allow_remote_admin')

# New enhanced options
enable_ipps=$(bashio::config 'enable_ipps')
ssl_cert=$(bashio::config 'ssl_cert')
ssl_key=$(bashio::config 'ssl_key')
require_encryption=$(bashio::config 'require_encryption')
max_jobs=$(bashio::config 'max_jobs')
max_jobs_per_printer=$(bashio::config 'max_jobs_per_printer')
job_history=$(bashio::config 'job_history')
preserve_job_history=$(bashio::config 'preserve_job_history')
preserve_job_files=$(bashio::config 'preserve_job_files')
access_log=$(bashio::config 'access_log')
page_log=$(bashio::config 'page_log')
max_log_size=$(bashio::config 'max_log_size')
debug_logging=$(bashio::config 'debug_logging')
raw_printing=$(bashio::config 'raw_printing')
auto_backup=$(bashio::config 'auto_backup')
cleanup_temp_files=$(bashio::config 'cleanup_temp_files')

# Set log level
case "${log_level}" in
    debug) LOG_LEVEL="debug2" ;;
    info) LOG_LEVEL="info" ;;
    warning) LOG_LEVEL="warn" ;;
    error) LOG_LEVEL="error" ;;
    *) LOG_LEVEL="info" ;;
esac

# Override with debug if enabled
if bashio::config.true 'debug_logging'; then
    LOG_LEVEL="debug2"
fi

# Set max log size
case "${max_log_size}" in
    1m) MAX_LOG_SIZE="1m" ;;
    5m) MAX_LOG_SIZE="5m" ;;
    10m) MAX_LOG_SIZE="10m" ;;
    50m) MAX_LOG_SIZE="50m" ;;
    100m) MAX_LOG_SIZE="100m" ;;
    unlimited) MAX_LOG_SIZE="0" ;;
    *) MAX_LOG_SIZE="10m" ;;
esac

# Set job preservation time
case "${preserve_job_files}" in
    1h) PRESERVE_TIME="3600" ;;
    6h) PRESERVE_TIME="21600" ;;
    1d) PRESERVE_TIME="86400" ;;
    7d) PRESERVE_TIME="604800" ;;
    30d) PRESERVE_TIME="2592000" ;;
    forever) PRESERVE_TIME="0" ;;
    *) PRESERVE_TIME="86400" ;;
esac

bashio::log.info "Setting up enhanced CUPS configuration..."
bashio::log.info "CUPS Version: $(cupsd --version 2>&1 || echo 'Version check failed')"

# Create enhanced cupsd.conf
cat > /etc/cups/cupsd.conf <<EOFCONF
# ==============================================================================
# CUPS 3 Enhanced Configuration
# Generated by Home Assistant Add-on
# Built from OpenPrinting CUPS source
# ==============================================================================

# Logging Configuration
LogLevel ${LOG_LEVEL}
MaxLogSize ${MAX_LOG_SIZE}
ErrorLog stderr
$(bashio::var.true "${access_log}" && echo "AccessLog stderr" || echo "AccessLog /dev/null")
$(bashio::var.true "${page_log}" && echo "PageLog stderr" || echo "PageLog /dev/null")

# Network Configuration
Listen 0.0.0.0:631
ServerAlias *
ServerName $(hostname)
ServerTokens Minimal

# SSL/TLS Configuration (IPPS)
$(if bashio::config.true 'enable_ipps'; then
    if [ -f "$ssl_cert" ] && [ -f "$ssl_key" ]; then
        echo "ServerCertificate $ssl_cert"
        echo "ServerKey $ssl_key"
        if bashio::config.true 'require_encryption'; then
            echo "DefaultEncryption Required"
        else
            echo "DefaultEncryption IfRequested"
        fi
    else
        echo "# IPPS enabled but SSL certificates not found"
        echo "DefaultEncryption Never"
    fi
else
    echo "DefaultEncryption Never"
fi)

# Browsing and Sharing
Browsing $(bashio::var.true "${browsing}" && echo "On" || echo "Off")
BrowseLocalProtocols dnssd
BrowseWebIF Yes
DefaultShared $(bashio::var.true "${share_printers}" && echo "Yes" || echo "No")

# Job Management
MaxJobs ${max_jobs}
MaxJobsPerPrinter ${max_jobs_per_printer}
MaxJobsPerUser 0
$(bashio::var.true "${job_history}" && echo "JobHistory On" || echo "JobHistory Off")
$(bashio::var.true "${preserve_job_history}" && echo "PreserveJobHistory Yes" || echo "PreserveJobHistory No")
PreserveJobFiles ${PRESERVE_TIME}

# Performance Options
MaxCopies 9999
MultipleOperationTimeout 300
Timeout 300

# Raw Printing
$(bashio::var.true "${raw_printing}" && echo "FileDevice Yes" || echo "FileDevice No")
$(bashio::var.true "${raw_printing}" && echo "RawPrint On" || echo "RawPrint Off")

# Security - Network Access
$(if bashio::config.has_value 'allowed_networks'; then
    echo "# Restricted Network Access"
    echo "<Location />"
    echo "  Order deny,allow"
    echo "  Deny from all"
    bashio::config 'allowed_networks' | jq -r '.[]' | while read -r network; do
        echo "  Allow from $network"
    done
    echo "</Location>"
else
    echo "# Allow all access"
    echo "<Location />"
    echo "  Order allow,deny"
    echo "  $(bashio::var.true "${allow_remote_admin}" && echo "Allow all" || echo "Allow @LOCAL")"
    echo "</Location>"
fi)

# Security - Admin Access
<Location /admin>
  Order allow,deny
  $(bashio::var.true "${allow_remote_admin}" && echo "Allow all" || echo "Allow @LOCAL")
</Location>

<Location /admin/conf>
  AuthType Default
  Require user @SYSTEM
  Order allow,deny
  $(bashio::var.true "${allow_remote_admin}" && echo "Allow all" || echo "Allow @LOCAL")
</Location>

# IPP Everywhere support
<Policy default>
  <Limit CUPS-Add-Modify-Printer CUPS-Delete-Printer CUPS-Add-Modify-Class CUPS-Delete-Class CUPS-Set-Default CUPS-Get-Devices>
    AuthType Default
    Require user @SYSTEM
    Order deny,allow
  </Limit>

  <Limit Pause-Printer Resume-Printer Enable-Printer Disable-Printer Pause-Printer-After-Current-Job Hold-New-Jobs Release-Held-New-Jobs Deactivate-Printer Activate-Printer Restart-Printer Shutdown-Printer Startup-Printer Promote-Job Schedule-Job-After Cancel-Jobs CUPS-Accept-Jobs CUPS-Reject-Jobs>
    AuthType Default
    Require user @SYSTEM
    Order deny,allow
  </Limit>

  <Limit CUPS-Authenticate-Job CUPS-Get-Document>
    Require user @OWNER @SYSTEM
    Order deny,allow
  </Limit>

  <Limit All>
    Order deny,allow
  </Limit>
</Policy>

# IPP Everywhere policy for driverless printing
<Policy driverless>
  <Limit All>
    Order deny,allow
  </Limit>
</Policy>
EOFCONF

# Create admin user if it doesn't exist
if ! id "${admin_user}" &>/dev/null; then
    bashio::log.info "Creating CUPS admin user: ${admin_user}"
    adduser -D -H -s /sbin/nologin "${admin_user}"
    echo "${admin_user}:${admin_password}" | chpasswd
    addgroup lpadmin 2>/dev/null || true
    addgroup "${admin_user}" lpadmin 2>/dev/null || true
fi

# Ensure directories exist
mkdir -p /var/cache/cups
mkdir -p /var/spool/cups/tmp
mkdir -p /var/run/cups
mkdir -p /config/cups-backups

# Set permissions
chown -R root:lpadmin /etc/cups 2>/dev/null || chown -R root:root /etc/cups
chmod 755 /etc/cups
chmod 640 /etc/cups/cupsd.conf

# Cleanup temporary files if enabled
if bashio::config.true 'cleanup_temp_files'; then
    bashio::log.debug "Cleaning up temporary files..."
    find /var/spool/cups/tmp -type f -mtime +1 -delete 2>/dev/null || true
fi

# Run initial backup if enabled
if bashio::config.true 'auto_backup'; then
    bashio::log.info "Running initial backup..."
    /usr/local/bin/cups-backup auto &
fi

# Configuration summary
bashio::log.info "=== CUPS 3 Enhanced Configuration ==="
bashio::log.info "Admin User: ${admin_user}"
bashio::log.info "WebUI: http://homeassistant.local:631"
bashio::log.info "Log Level: ${log_level}"
bashio::log.info "Share Printers: ${share_printers}"
bashio::log.info "Browsing: ${browsing}"
bashio::log.info "IPPS: ${enable_ipps}"
bashio::log.info "Max Jobs: ${max_jobs}"
bashio::log.info "Auto Backup: ${auto_backup}"
bashio::log.info "Auto Discovery: $(bashio::config 'auto_discovery')"
bashio::log.info "======================================"

# Start CUPS in foreground
exec /usr/sbin/cupsd -f
