# rootfs/etc/services.d/avahi/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Handle Avahi daemon shutdown
# ==============================================================================
bashio::log.info "Avahi daemon stopped"
